# docker/openresty.ini
[program:openresty]
command = /usr/local/openresty/nginx/sbin/nginx
startsecs = 0
stopwaitsecs = 0
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /opt/opsant/logs/openresty.log
