# docker/opsant.ini
[program: opsant]
command = gunicorn wsgi --bind 127.0.0.1:8001 -k gevent -w 2 -n opsant-web --access-logfile - --error-logfile -
directory = /opt/opsant-backend
environment = DJANGO_SETTINGS_MODULE=settings,OPS_ANT_ENV="production",PAAS_LOGGING_DIR="/opt/opsant/logs"
startsecs = 0
stopwaitsecs = 0
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /opt/opsant/logs/opsant.log


[program: opsant_celery]
command = python3 /opt/opsant-backend/manage.py celery worker -n opsant-celery -l INFO --concurrency=2
directory = /opt/opsant-backend
environment = DJANGO_SETTINGS_MODULE=settings,OPS_ANT_ENV="production",PAAS_LOGGING_DIR="/opt/opsant/logs"
startsecs = 0
stopwaitsecs = 0
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /opt/opsant/logs/opsant_celery.log


[program: opsant_beat]
command = python3 /opt/opsant-backend/manage.py celery beat
directory = /opt/opsant-backend
environment = DJANGO_SETTINGS_MODULE=settings,OPS_ANT_ENV="production",PAAS_LOGGING_DIR="/opt/opsant/logs"
startsecs = 0
stopwaitsecs = 0
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /opt/opsant/logs/opsant_celery.log
