# docker/websocket.ini
[program:websocket]
command = daphne --proxy-headers -b 127.0.0.1 -p 8002 asgi:application
directory = /opt/opsant-backend
environment = OPS_ANT_ENV="production",BK_LOG_DIR="/opt/opsant/logs"
startsecs = 0
stopwaitsecs = 0
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /opt/opsant/logs/websocket.log
